# fast-oom/fast-oom.service
[Unit]
Description=Trigger OOM Killer fast

[Service]
Type=exec
ExecStart=/usr/local/sbin/fast-oom

[Install]
WantedBy=default.target
